# host-elementary/config/autostart/disable_capslock.desktop
[Desktop Entry]
Name=Disable Capslock
Exec=setxkbmap -option caps:escape
Terminal=false
Type=Application
