# host-elementary/config/autostart/keyboard-repeat.desktop
[Desktop Entry]
Name=Keyboard Repeat
Exec=xset r rate 150 75
Terminal=false
Type=Application
